# Small_5/package/ipk/luci-app-vsftpd/files/etc/uci-defaults/luci-vsftpd
uci batch <<EOF
	delete ucitrack.@vsftpd[-1]
	add ucitrack vsftpd
	set ucitrack.@vsftpd[-1].init=vsftpd
	commit ucitrack
EOF
