# Small_5/package/ipk/luci-app-vlmcsd/root/etc/uci-defaults/luci-vlmcsd
uci batch <<EOF
	delete ucitrack.@vlmcsd[-1]
	add ucitrack vlmcsd
	set ucitrack.@vlmcsd[-1].init=vlmcsd
	commit ucitrack
EOF
