# Small_5/package/ipk/luci-app-openvpn-server/root/etc/uci-defaults/luci-openvpn
[ -s /etc/openvpn/ta.key ] || openvpn --genkey --secret /etc/openvpn/ta.key
