# Small_5/package/ipk/luci-app-filetransfer/etc/uci-defaults/luci-filetransfer
sed -i 's/cbi.submit\"] = true/cbi.submit\"] = \"1\"/g' /usr/lib/lua/luci/dispatcher.lua
exit 0
