# Small_5/package/ipk/luci-app-cpufreq/root/etc/uci-defaults/cpufreq
uci batch <<EOF
	delete ucitrack.@cpufreq[-1]
	add ucitrack cpufreq
	set ucitrack.@cpufreq[-1].init=cpufreq
	commit ucitrack
EOF
